# config/etc/supervisor/conf.d/redis.conf
[program:redis]
command=/usr/bin/redis-server /etc/redis/redis.conf
redirect_stderr=true
stopsignal=TERM
autorestart=true
